# custom-files/tl-wr840n-v5/etc/uci-defaults/z001_flashbox-system.sh
../../../tl-wr840n-v4/etc/uci-defaults/z001_flashbox-system.sh