# files/etc/init.d/anlixboot
#!/bin/sh /etc/rc.common

. /usr/share/libubox/jshn.sh
if [ -e /usr/share/functions/custom_device.sh ]
then
	. /usr/share/functions/network_functions.sh
fi

START=17

boot() {
	[ "$(type -t anlix_bootup_defaults)" ] && anlix_bootup_defaults

	if [ "$(type -t set_switch_bridge_mode_on_boot)" ]
	then
		update_vlan "n"
	fi
}
