# files/etc/init.d/minisapo
#!/bin/sh /etc/rc.common

. /usr/share/functions/device_functions.sh
. /usr/share/functions/network_functions.sh
. /usr/share/functions/wireless_functions.sh

START=9F
STOP=1F

USE_PROCD=1

start_service() {
  BRIDGE_IP_ADDR="$(get_lan_bridge_ipaddr)"
  MAC="$(get_mac)"
  SSID="$(get_wifi_local_config | jsonfilter -e '@["local_ssid_24"]')"
  SSID_5="$(get_wifi_local_config | jsonfilter -e '@["local_ssid_50"]')"

  procd_open_instance
  procd_set_param command minisapo "$BRIDGE_IP_ADDR $MAC $SSID $SSID_5"
  procd_set_param respawn
  procd_set_param stdout 1 # forward stdout of the command to logd
  procd_set_param stderr 1 # same for stderr
  procd_set_param pidfile /var/run/minisapo.pid
  procd_close_instance
}
