# custom-files/tl-wr849n-v6/etc/uci-defaults/z005_flashbox-wireless.sh
../../../tl-wr840n-v4/etc/uci-defaults/z005_flashbox-wireless.sh