# custom-files/dir-815-d1/etc/init.d/rtl8192cd
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=18

boot() {
	ifconfig wlan0 up
	ifconfig wlan1 up

	iwpriv wlan0 set_mib xcap=30
	iwpriv wlan0 set_mib ther=32
	iwpriv wlan0 set_mib pwrlevelCCK_A=2828282828282828282828282828
	iwpriv wlan0 set_mib pwrlevelCCK_B=2d2d2d2d2d2d2d2d2d2d2d2d2d2d
	iwpriv wlan0 set_mib pwrlevelHT40_1S_A=2c2c2c2c2c2c2c2c2c2c2c2c2c2c
	iwpriv wlan0 set_mib pwrlevelHT40_1S_B=3030303030303030303030303030
	
	iwpriv wlan1 set_mib xcap=40
	iwpriv wlan1 set_mib ther=27
	iwpriv wlan1 set_mib pwrlevel5GHT40_1S_A=0000000000000000000000000000000000000000000000000000000000000000000000262626262626262626262626262020202020202020202020202020202021212121212121212121212121212121212121212121212121212121212121212121212121212121212121212121212121212121212121211d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d1d00000000000000000000000000000000000000
	iwpriv wlan1 set_mib pwrlevel5GHT40_1S_B=0000000000000000000000000000000000000000000000000000000000000000000000262626262626262626262626262020202020202020202020202020202022222222222222222222222222222222222222222222222222222222222222222222222222222222222222222222222222222222222222221b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b1b00000000000000000000000000000000000000
	iwpriv wlan1 set_mib pwrdiff_5G_80BW1S_160BW1S_A=0000000000000000000000000000000000000000000000000000000000000000000000d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d000000000000000000000000000000000000000
	iwpriv wlan1 set_mib pwrdiff_5G_80BW1S_160BW1S_B=0000000000000000000000000000000000000000000000000000000000000000000000d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d0d000000000000000000000000000000000000000

	ifconfig wlan0 down
	ifconfig wlan1 down
}
